# crontab
* * * * * D:/wamp64/bin/php/php7.0.33/php.exe D:/wamp64/www/rrr/New folder/Tech2Success/Index.php